# infra/aws/systemd/mousekeeper-postgres-backup.timer
[Unit]
Description=Run MOUSEKEEPER PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 03:20:00 UTC
RandomizedDelaySec=30m
Persistent=true
Unit=mousekeeper-postgres-backup.service

[Install]
WantedBy=timers.target
